# etc/cron.d/picam
MAILTO=""
PATH=/bin:/usr/bin:/usr/local/bin
*   *  *   *   * motion /usr/local/bin/picam_supervise control 2>&1 | logger -t picam_supervise


